# root/etc/services.d/unionfs/run
#!/usr/bin/with-contenv sh
unionfs -o cow,allow_other,auto_cache,nonempty /source1=RW:/source2=RO /target
tail -f /dev/null
